# nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off; master_process on;'
stdout_events_enabled=true
stderr_events_enabled=true