# phpfpm.conf
[program:phpfpm]
command=/usr/sbin/php-fpm7.0 --nodaemonize -c /etc/php/7.0/fpm --fpm-config /etc/php/7.0/fpm/php-fpm.conf
stdout_events_enabled=true
stderr_events_enabled=true