# phpfpm.conf
[program:phpfpm]
command=/usr/sbin/php-fpm7.0 -c /etc/php/7.0/fpm --fpm-config /etc/php/7.0/fpm/php-fpm.conf